# launchagents/com.user.media-stack.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.media-stack</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>__SCRIPT_DIR__/scripts/start_stack.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/media-stack.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/media-stack-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
